# packaging/systemd/wyrelog-system.service
[Unit]
Description=Wyrelog system profile service
Documentation=man:wyrelogd(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wyrelog
Group=wyrelog
Environment=WYL_LOG=warn
EnvironmentFile=-/etc/wyrelog/system.env
LoadCredential=wyrelog-system-policy-key:/etc/wyrelog/system/policy.key
ExecStart=/usr/bin/wyrelogd --config /etc/wyrelog/wyrelogd.conf --production
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/etc/wyrelog /usr/share/wyrelog
# Defense-in-depth: the line above already covers /etc/wyrelog/, but
# pinning the conf file explicitly survives a future drop-in that
# relaxes the directory rule. The daemon also gates on mode/owner at
# open(2) (wyl_daemon_options_resolve in wyrelog/daemon/options.c).
ReadOnlyPaths=/etc/wyrelog/wyrelogd.conf
ReadWritePaths=/var/lib/wyrelog/system /var/log/wyrelog/system /run/wyrelog
RuntimeDirectory=wyrelog
RuntimeDirectoryMode=0750
StateDirectory=wyrelog
StateDirectoryMode=0750
LogsDirectory=wyrelog
LogsDirectoryMode=0750
UMask=0077

[Install]
WantedBy=multi-user.target
